# dockerfiles/crontab
* * * * * cd /var/www/html/project && php artisan dispatch:job 1 && php artisan dispatch:job 2 && php artisan dispatch:job 3 && php artisan schedule:run >> /dev/stdout 2>&1
