# dockerfiles/crontab
* * * * * cd /var/www/html/project && php artisan schedule:run >> /dev/stdout 2>&1
